# project-templates/obsidian-vault/launchd/com.motorthings.granola-enricher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.motorthings.granola-enricher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/motorthings/Documents/Vault/_tools/scripts/granola_enricher.py</string>
    </array>

    <!-- Fire when anything changes in the Granola folder -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/motorthings/Documents/Vault/Granola/Summaries</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/motorthings/Library/Logs/granola-enricher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/motorthings/Library/Logs/granola-enricher-error.log</string>
</dict>
</plist>
